# systemd/.config/systemd/user/input-remapper-watchdog.timer
[Unit]
Description=Input Remapper Watchdog Timer
After=graphical-session.target

[Timer]
OnBootSec=30s
OnUnitActiveSec=30s
AccuracySec=5s

[Install]
WantedBy=timers.target
