# systemd/.config/systemd/user/swayidle.service
[Unit]
Description=Swayidle - Idle management daemon for Wayland
Documentation=man:swayidle(1)
PartOf=graphical-session.target
After=graphical-session.target

[Service]
Type=simple
ExecStart=/run/current-system/sw/bin/swayidle -w timeout 570 'notify-send -u critical -t 25000 -i %h/.config/waybar/caution_4539472.png "🔒 Verrouillage imminent" "L՚écran se verrouillera dans 30 secondes"' resume 'swaync-client -C' timeout 600 '%h/.config/waybar/safe-lock.sh' timeout 900 'systemctl suspend' before-sleep '%h/.config/waybar/safe-lock.sh' after-resume '%h/nixos-config/bin/resume-actions.sh'
Restart=always
RestartSec=5
KillMode=mixed
TimeoutStopSec=10

# Restart periodically to prevent stale state (every 6 hours)
RuntimeMaxSec=6h

[Install]
WantedBy=graphical-session.target
